# root/etc/service/apache2/run
#!/bin/sh
exec 2>&1

exec /usr/sbin/httpd -D FOREGROUND -f /web/config/httpd.conf

